# infra/quadlet/centsible-mariadb.container
# Centsible MariaDB Container (runs in centsible pod)
# Database port 3306 is NOT exposed externally.
# Only accessible from other containers in the same pod via localhost.

[Unit]
Description=Centsible MariaDB Database (Pod)
Documentation=https://mariadb.org/documentation/
After=network-online.target
BindsTo=centsible-pod.service
After=centsible-pod.service

[Container]
Image=docker.io/library/mariadb:11.7
ContainerName=centsible-mariadb
AutoUpdate=registry

# Run in the centsible pod
Pod=centsible.pod

# Persistent data volume
Volume=centsible-db:/var/lib/mysql:Z

# Load environment variables (MARIADB_ROOT_PASSWORD, MARIADB_DATABASE, etc.)
EnvironmentFile=%h/.config/containers/systemd/.env.centsible

# Health check
HealthCmd=healthcheck.sh --connect --innodb_initialized
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=default.target
