# infra/quadlet/centsible-dev.container
# Centsible Dev Utility Container (runs in centsible pod)
# Utility-only: verification, exec-based workflows, and maintenance commands.
# No app ports are exposed from this container.

[Unit]
Description=Centsible Dev Utility Container (Pod)
Documentation=https://podman.io/
After=network-online.target centsible-pod.service
Wants=network-online.target
BindsTo=centsible-pod.service

[Container]
Image=localhost/centsible-dev:latest
ContainerName=centsible-dev
AutoUpdate=local
Pod=centsible.pod
WorkingDir=/workspace
Exec=sleep infinity
Environment=NEXT_TELEMETRY_DISABLED=1

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
